# infra/hermes/systemd/bairui-hermes.service
[Unit]
Description=bairui Hermes Agent OS runtime
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/bairui/hermes
EnvironmentFile=/etc/bairui/hermes.env
ExecStart=/usr/bin/python -m src.hermes
Restart=always
RestartSec=5
User=bairui
Group=bairui
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/bairui /var/log/bairui /opt/bairui/hermes

[Install]
WantedBy=multi-user.target
